# deploy/serverbee-agent.service
[Unit]
Description=ServerBee Agent
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/serverbee-agent
WorkingDirectory=/etc/serverbee
Restart=always
RestartSec=5
# Exit code 78 = permanent enrollment-code failure; don't restart, the operator
# must rotate the code (otherwise we burn the registration rate limit).
RestartPreventExitStatus=78
AmbientCapabilities=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
